# openvpn.conf
start on desktop-start
stop on desktop-end

script
  cd ~/apps/systray-openvpn/ && nohup python test.py &
end script
